# crontab
8 * * * * /home/opc/waste.sh > /home/opc/waste.log 2>&1
